# com.ai-news.daily-poster.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-news.daily-poster</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wzh/project/news/scripts/daily-ai-poster-plus_claude.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/wzh/project/news</string>

    <key>StandardOutPath</key>
    <string>/Users/wzh/project/news/ai_news/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wzh/project/news/ai_news/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/wzh</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
